# manager/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

# ============================================================
# Infrastructure Layer (start first)
# ============================================================

[program:minio]
command=/opt/hiclaw/scripts/init/start-minio.sh
priority=50
autorestart=true
stdout_logfile=/var/log/hiclaw/minio.log
stderr_logfile=/var/log/hiclaw/minio-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:tuwunel]
command=/opt/hiclaw/scripts/init/start-tuwunel.sh
priority=100
autorestart=true
stdout_logfile=/var/log/hiclaw/tuwunel.log
stderr_logfile=/var/log/hiclaw/tuwunel-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# ============================================================
# Element Web (Nginx static file server)
# ============================================================

[program:element-web]
command=/opt/hiclaw/scripts/init/start-element-web.sh
priority=650
autorestart=true
stdout_logfile=/var/log/hiclaw/element-web.log
stderr_logfile=/var/log/hiclaw/element-web-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# ============================================================
# Application Layer
# ============================================================

[program:setup-host-symlinks]
command=/opt/hiclaw/scripts/setup-host-symlinks.sh
priority=750
autostart=true
autorestart=false
startretries=0
startsecs=0
stdout_logfile=/var/log/hiclaw/setup-host-symlinks.log
stderr_logfile=/var/log/hiclaw/setup-host-symlinks-error.log
stdout_logfile_maxbytes=1MB

[program:mc-mirror]
command=/opt/hiclaw/scripts/init/start-mc-mirror.sh
priority=700
autorestart=true
stdout_logfile=/var/log/hiclaw/mc-mirror.log
stderr_logfile=/var/log/hiclaw/mc-mirror-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:manager-agent]
command=/opt/hiclaw/scripts/init/start-manager-agent.sh
priority=800
autorestart=true
stdout_logfile=/var/log/hiclaw/manager-agent.log
stderr_logfile=/var/log/hiclaw/manager-agent-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
